# osmanager/init/luoyun/customconfig.conf
start on stopped luoyun/osconfig

console output
task
script
  if [ -x /LuoYun/custom/scripts/config ]
  then
    /LuoYun/custom/scripts/config 
  else
    exit 0
  fi
end script
