# osmanager/init/luoyun/customweb.conf
start on stopped luoyun/osconfig

console output
respawn
exec /LuoYun/bin/pyweb/lyweb.py 8080
