# osmanager/init/luoyun/osconfig.conf
start on stopped rc RUNLEVEL=[345]

console output
task
exec /LuoYun/bin/lyconfig.sh
