# osmanager/init/luoyun/osm.conf
start on stopped luoyun/osconfig

console output
respawn
exec /LuoYun/bin/pyosm/lyosm.py
