# tools/init/nginx.conf
start on started rc RUNLEVEL=[345]

console output
expect daemon
respawn
exec nginx
